# supervisord.conf
[supervisord]
nodaemon=true

[program:fdbserver]
command=/bin/bash /usr/lib/foundationdb/docker-start.sh

[program:fdbsql]
command=/bin/bash /usr/lib/foundationdb/docker-sql-layer.sh
